# my_portfolio_addon/rootfs/etc/services.d/my_service/run
chmod +x rootfs/etc/services.d/my_service/run
